# custom-upload.service
[Unit]
Description=Custom Upload for Products Design - Shopify App
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/3d-customizer
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/pnpm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=custom-upload

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/3d-customizer

[Install]
WantedBy=multi-user.target

